# madarles/systemd/smbd.service.d/drop.conf
[Unit]
BindsTo=mnt-ADATA.mount
After=mnt-ADATA.mount
Requires=mnt-ADATA.mount
PartOf=mnt-ADATA.mount

[Service]
Restart=always
RestartSec=20s

[Install]
WantedBy=mnt-ADATA.mount
