# madarles/systemd/emergency.timer
[Unit]
Description=Emergency timer

[Timer]
OnUnitActiveSec=45s
Unit=status.service
