# madarles/systemd/status.service
[Unit]
Description=Status service
After=multi-user.target

[Service]
Type=forking
ExecStart=/home/pi/bin/check_status.sh
KillMode=process

