# madarles/systemd/status.service
[Unit]
Description=Status service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/home/pi/bin/check_status.sh

