# madarles/systemd/status.timer
[Unit]
Description=Status timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=status.service

[Install]
WantedBy=multi-user.target
