# madarles/systemd/upload.service
[Unit]
Description=Upload service
After=multi-user.target

[Service]
Type=oneshot
User=pi
ExecStart=/home/pi/bin/upload

