# madarles/systemd/upload.timer
[Unit]
Description=Upload timer

[Timer]
OnBootSec=10min
OnCalendar=daily
Unit=upload.service

[Install]
WantedBy=multi-user.target
